# files/etc/systemd/system/aide-fim.service
[Unit]
Description=AIDE File Integrity Monitoring Service

[Service]
Type=simple
ExecStart=/usr/bin/aide --config /etc/aide/aide.conf --check
